# dbbackup.service
[Unit]
Description=DB Backup

[Service]
ExecStart=/bin/bash /path/db_bak.sh